# .ansible/playbooks/archon/roles/wallpaper-rng/templates/wallpaper-rng.service.j2
# {{ ansible_managed }}

[Unit]
Description=Randomly change the current wallpaper

[Service]
ExecStart={{ ansible_env.HOME }}/.local/scripts/rofi_wallpapers.sh --rng --silent
Type=oneshot